# debian/minios-apparmor.service
[Unit]
Description=Changes apparmor service

[Service]
Type=oneshot
ExecStart=/usr/sbin/minios-apparmor-patch
TimeoutStartSec=60

[Install]
RequiredBy=minios-apparmor.path
